# service/repo-hooker.service
[Unit]
Description=Pulls repos in automatigically

[Service]
Type=simple
ExecStart=bash /usr/bin/repo-hooker.sh

[Install]
WantedBy=multi-user.target